# usr/lib/systemd/user/pulseaudio.service.d/10-rosemary.conf
[Service]
Environment=HYBRIS_LD_LIBRARY_PATH="/android/system/lib64:/android/vendor/lib64:/android/system/apex/com.android.art.release/lib64/"
